# .config/systemd/user/i3lock.service
[Unit]
Description=i3lock: An improved screenlocker based upon XCB and PAM

[Service]
Type=forking
ExecStart=/home/muldvang/.local/bin/i3lock-blur
Environment=DISPLAY=:0

[Install]
WantedBy=default.target
